# crontab
10 */6 * * * /usr/local/bin/php /var/www/symfony/app/console api:load:rates
