# crontab
10 */6 * * * php /var/www/symfony/app/console api:load:rates
